# Deployment/var-origin/lib/dpkg/info/fwupd.list
/.
/etc
/etc/fwupd
/etc/fwupd/bios-settings.d
/etc/fwupd/bios-settings.d/README.md
/etc/fwupd/fwupd.conf
/etc/fwupd/remotes.d
/etc/fwupd/remotes.d/lvfs-testing.conf
/etc/fwupd/remotes.d/lvfs.conf
/etc/fwupd/remotes.d/vendor-directory.conf
/etc/grub.d
/etc/grub.d/35_fwupd
/etc/pki
/etc/pki/fwupd
/etc/pki/fwupd/GPG-KEY-Linux-Foundation-Firmware
/etc/pki/fwupd/GPG-KEY-Linux-Vendor-Firmware-Service
/etc/pki/fwupd/LVFS-CA.pem
/etc/pki/fwupd-metadata
/etc/pki/fwupd-metadata/GPG-KEY-Linux-Foundation-Metadata
/etc/pki/fwupd-metadata/GPG-KEY-Linux-Vendor-Firmware-Service
/etc/pki/fwupd-metadata/LVFS-CA.pem
/etc/update-motd.d
/etc/update-motd.d/85-fwupd
/usr
/usr/bin
/usr/bin/dbxtool
/usr/bin/fwupdmgr
/usr/bin/fwupdtool
/usr/lib
/usr/lib/modules-load.d
/usr/lib/modules-load.d/fwupd-msr.conf
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/systemd/system/fwupd-offline-update.service
/usr/lib/systemd/system/fwupd-refresh.service
/usr/lib/systemd/system/fwupd-refresh.timer
/usr/lib/systemd/system/fwupd.service
/usr/lib/systemd/system/system-update.target.wants
/usr/lib/systemd/system-shutdown
/usr/lib/systemd/system-shutdown/fwupd.shutdown
/usr/lib/sysusers.d
/usr/lib/sysusers.d/fwupd.conf
/usr/lib/udev
/usr/lib/udev/rules.d
/usr/lib/udev/rules.d/90-fwupd-devices.rules
/usr/lib/x86_64-linux-gnu
/usr/lib/x86_64-linux-gnu/fwupd-1.9.30
/usr/lib/x86_64-linux-gnu/fwupd-1.9.30/libfu_plugin_flashrom.so
/usr/lib/x86_64-linux-gnu/fwupd-1.9.30/libfu_plugin_modem_manager.so
/usr/lib/x86_64-linux-gnu/fwupd-1.9.30/libfwupdengine.so
/usr/lib/x86_64-linux-gnu/fwupd-1.9.30/libfwupdplugin.so
/usr/lib/x86_64-linux-gnu/fwupd-1.9.30/libfwupdutil.so
/usr/libexec
/usr/libexec/fwupd
/usr/libexec/fwupd/fwupd
/usr/libexec/fwupd/fwupd-detect-cet
/usr/libexec/fwupd/fwupdoffline
/usr/share
/usr/share/bash-completion
/usr/share/bash-completion/completions
/usr/share/bash-completion/completions/fwupdmgr
/usr/share/bash-completion/completions/fwupdtool
/usr/share/dbus-1
/usr/share/dbus-1/interfaces
/usr/share/dbus-1/interfaces/org.freedesktop.fwupd.xml
/usr/share/dbus-1/system-services
/usr/share/dbus-1/system-services/org.freedesktop.fwupd.service
/usr/share/dbus-1/system.d
/usr/share/dbus-1/system.d/org.freedesktop.fwupd.conf
/usr/share/doc
/usr/share/doc/fwupd
/usr/share/doc/fwupd/copyright
/usr/share/fish
/usr/share/fish/vendor_completions.d
/usr/share/fish/vendor_completions.d/fwupdmgr.fish
/usr/share/fwupd
/usr/share/fwupd/add_capsule_header.py
/usr/share/fwupd/firmware_packager.py
/usr/share/fwupd/host-emulate.d
/usr/share/fwupd/host-emulate.d/thinkpad-p1-iommu.json.gz
/usr/share/fwupd/install_dell_bios_exe.py
/usr/share/fwupd/metainfo
/usr/share/fwupd/metainfo/org.freedesktop.fwupd.remotes.lvfs-testing.metainfo.xml
/usr/share/fwupd/metainfo/org.freedesktop.fwupd.remotes.lvfs.metainfo.xml
/usr/share/fwupd/quirks.d
/usr/share/fwupd/quirks.d/builtin.quirk.gz
/usr/share/fwupd/remotes.d
/usr/share/fwupd/remotes.d/vendor
/usr/share/fwupd/remotes.d/vendor/firmware
/usr/share/fwupd/remotes.d/vendor/firmware/README.md
/usr/share/fwupd/simple_client.py
/usr/share/fwupd/uefi-capsule-ux.tar.xz
/usr/share/icons
/usr/share/icons/hicolor
/usr/share/icons/hicolor/scalable
/usr/share/icons/hicolor/scalable/apps
/usr/share/icons/hicolor/scalable/apps/org.freedesktop.fwupd.svg
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/fwupd
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/dbxtool.1.gz
/usr/share/man/man1/fwupdmgr.1.gz
/usr/share/man/man1/fwupdtool.1.gz
/usr/share/man/man5
/usr/share/man/man5/fwupd-remotes.d.5.gz
/usr/share/man/man5/fwupd.conf.5.gz
/usr/share/man/man8
/usr/share/man/man8/fwupd-refresh.service.8.gz
/usr/share/metainfo
/usr/share/metainfo/org.freedesktop.fwupd.metainfo.xml
/usr/share/polkit-1
/usr/share/polkit-1/actions
/usr/share/polkit-1/actions/org.freedesktop.fwupd.policy
/usr/share/polkit-1/rules.d
/usr/share/polkit-1/rules.d/org.freedesktop.fwupd.rules
/var
/var/cache
/var/cache/app-info
/var/cache/app-info/xmls
/var/lib
/var/lib/polkit-1
/var/lib/polkit-1/localauthority
/var/lib/polkit-1/localauthority/10-vendor.d
/var/lib/polkit-1/localauthority/10-vendor.d/fwupd.pkla
/usr/lib/systemd/system/system-update.target.wants/fwupd-offline-update.service
/usr/share/doc/fwupd/changelog.Debian.gz
