# Deployment/var-origin/lib/dpkg/info/screen.postinst
#!/bin/sh
set -e

. /usr/share/debconf/confmodule

if [ "$1" = configure ]; then
  if ! test -d /run/screen; then
    install -g utmp -m 0775 -d /run/screen
  fi
  perms="`stat -c%a /usr/bin/screen`"
  override=/etc/tmpfiles.d/screen-cleanup.conf
  if [ $perms -eq 4755 ]; then
    chmod 0755 /run/screen
    if [ ! -f $override ]; then
       echo "# This file is generated by $0 upon package configuration" > $override
       echo 'd /run/screen 0755 root utmp' >> $override
    fi
  elif [ $perms -eq 755 ]; then
    chmod 1777 /run/screen
    if [ ! -f $override ]; then
       echo "# This file is generated by $0 upon package configuration" > $override
       echo 'd /run/screen 1777 root utmp' >> $override
    fi
  fi

  # LP #1462692
  servicedir=/lib/systemd/system
  servicefile=$servicedir/screen-cleanup.service
  mkdir -p $servicedir
  if [ ! -e $servicefile ]; then
      ln -s /dev/null $servicefile
  fi
fi

# Automatically added by dh_installtmpfiles/13.14.1ubuntu5
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -x "$(command -v systemd-tmpfiles)" ]; then
		systemd-tmpfiles ${DPKG_ROOT:+--root="$DPKG_ROOT"} --create screen-cleanup.conf || true
	fi
fi
# End automatically added section
# Automatically added by dh_installinit/13.14.1ubuntu5
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -x "/etc/init.d/screen-cleanup" ]; then
		update-rc.d screen-cleanup defaults >/dev/null || exit 1
	fi
fi
# End automatically added section

