# Deployment/var/lib/dpkg/info/certbot.list
/.
/etc
/etc/cron.d
/etc/cron.d/certbot
/etc/letsencrypt
/etc/letsencrypt/cli.ini
/etc/logrotate.d
/etc/logrotate.d/certbot
/usr
/usr/bin
/usr/bin/certbot
/usr/lib
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/systemd/system/certbot.service
/usr/lib/systemd/system/certbot.timer
/usr/share
/usr/share/doc
/usr/share/doc/certbot
/usr/share/doc/certbot/README.rst
/usr/share/doc/certbot/copyright
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/certbot
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/certbot.1.gz
/usr/share/man/man7
/usr/share/man/man7/certbot.7.gz
/usr/bin/letsencrypt
/usr/share/doc/certbot/changelog.Debian.gz
/usr/share/man/man1/letsencrypt.1.gz
