# Deployment/var/lib/dpkg/info/dmeventd.list
/.
/usr
/usr/lib
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/systemd/system/dm-event.service
/usr/lib/systemd/system/dm-event.socket
/usr/lib/x86_64-linux-gnu
/usr/lib/x86_64-linux-gnu/device-mapper
/usr/lib/x86_64-linux-gnu/device-mapper/libdevmapper-event-lvm2mirror.so
/usr/lib/x86_64-linux-gnu/device-mapper/libdevmapper-event-lvm2raid.so
/usr/lib/x86_64-linux-gnu/device-mapper/libdevmapper-event-lvm2snapshot.so
/usr/lib/x86_64-linux-gnu/device-mapper/libdevmapper-event-lvm2thin.so
/usr/lib/x86_64-linux-gnu/device-mapper/libdevmapper-event-lvm2vdo.so
/usr/lib/x86_64-linux-gnu/libdevmapper-event-lvm2.so.2.03
/usr/sbin
/usr/sbin/dmeventd
/usr/share
/usr/share/doc
/usr/share/doc/dmeventd
/usr/share/doc/dmeventd/copyright
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/dmeventd.8.gz
/usr/lib/x86_64-linux-gnu/libdevmapper-event-lvm2.so
/usr/lib/x86_64-linux-gnu/libdevmapper-event-lvm2mirror.so
/usr/lib/x86_64-linux-gnu/libdevmapper-event-lvm2raid.so
/usr/lib/x86_64-linux-gnu/libdevmapper-event-lvm2snapshot.so
/usr/lib/x86_64-linux-gnu/libdevmapper-event-lvm2thin.so
/usr/lib/x86_64-linux-gnu/libdevmapper-event-lvm2vdo.so
/usr/share/doc/dmeventd/changelog.Debian.gz
