# Deployment/var/lib/dpkg/info/pollinate.list
/.
/etc
/etc/default
/etc/default/pollinate
/etc/pollinate
/etc/pollinate/entropy.ubuntu.com.pem
/usr
/usr/bin
/usr/bin/pollinate
/usr/lib
/usr/lib/nagios
/usr/lib/nagios/plugins
/usr/lib/nagios/plugins/check_pollen
/usr/lib/systemd
/usr/lib/systemd/system
/usr/lib/systemd/system/pollinate.service
/usr/share
/usr/share/doc
/usr/share/doc/pollinate
/usr/share/doc/pollinate/changelog.Debian.gz
/usr/share/doc/pollinate/copyright
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/pollinate.1.gz
